# root/etc/services.d/tinyproxy/run
#!/usr/bin/with-contenv bash

tinyproxy -c /config/proxy.conf
