# root/etc/services.d/tinyproxy/run
#!/usr/bin/with-contenv bash

tinyproxy -dc /etc/tinyproxy/tinyproxy.conf
